# systemd/hello-monitor.timer
[Unit]
Description=Timer for hello-monitor.service

[Timer]
OnBootSec=10
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
